# ssh-proxy/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:fastapi]
command=uvicorn server.main:app --host 0.0.0.0 --port 8022
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/fastapi.err.log
stdout_logfile=/var/log/supervisor/fastapi.out.log
priority=10

[program:tftpd]
command=/usr/sbin/in.tftpd --foreground --listen --user root --address 0.0.0.0:69 --secure /opt/gns3-ssh-proxy/tftp --create --permissive
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/tftpd.err.log
stdout_logfile=/var/log/supervisor/tftpd.out.log
priority=20

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
priority=30
